# salt/roots/salt/roles/bccvl/visualiser/visualiser_supervisord.ini
{% set user = salt['pillar.get']('visualiser:user', {'name': 'visualiser'}) %}

[program:visualiser]
command = scl enable python27 ". bin/activate; /home/{{ user.name }}/BCCVL_Visualiser/BCCVL_Visualiser/bin/pserve bccvl_production.ini"
process_name = visualiser
directory = /home/{{ user.name }}/BCCVL_Visualiser/BCCVL_Visualiser
priority = 20
redirect_stderr = true
user = {{ user.name }}
stdout_logfile=syslog
stderr_logfile=syslog
autostart = true
numprocs = 1

; When resorting to send SIGKILL, we want to kill all child processes as well
killasgroup = true
; Make sure subprocesses are stopped along with main process
stopasgroup = true