# salt/roots/salt/roles/worker/worker_supervisor.ini
{% set user = salt['pillar.get']('worker:user', {'name': 'bccvl'}) %}
{% set data_mover = salt['pillar.get']('worker:data_mover', '192.168.100.200') %}

[program:worker]
command = scl enable python27 ". bin/activate; bin/celery worker --app=org.bccvl.tasks --loglevel=INFO  --queues=worker --no-color -Ofair --hostname=worker@%%h"
environment = CELERY_JSON_CONFIG="/home/{{ user.name }}/worker/celery.json",WORKDIR="/tmp",DATA_MOVER="https://{{ data_mover }}/_data_mover/data_mover",HOME="/home/{{ user.name }}",USER="{{ user.name }}",MAXENT="/home/{{ user.name }}/maxent.jar",R_LIBS="/home/{{ user.name }}/R_libs"
process_name = worker
directory = /home/{{ user.name }}/worker/
priority = 998
user = {{ user.name }}
stdout_logfile=syslog
stderr_logfile=syslog
autostart = true
numprocs = 1
startsecs = 10

; Wait for current tasks to finish at shutdown
; should be at least the time of the longest expected job
stopwaitsecs = 600

; When resorting to send SIGKILL, we want to kill all child processes as well
killasgroup = true
; Make sure subprocesses are stopped along with main process
stopasgroup = true