# salt/roots/salt/roles/bccvl/plone/plone_worker_supervisord.ini
{% set user = salt['pillar.get']('plone:user', {'name': 'plone'}) %}

[program:plone_worker]
command = scl enable python27 ". bin/activate; bin/celery worker --app=org.bccvl.tasks --loglevel=INFO  --queues=plone --no-color -I org.bccvl.tasks.plone --hostname=plone@%%h"
environment = BCCVL_CONFIG="/home/{{ user.name }}/bccvl_buildout/etc/bccvl.ini",Z_CONFIG_FILE="parts/instance1/etc/zope.conf",CELERY_CONFIG_MODULE="celeryconfig"
process_name = plone_worker
directory = /home/{{ user.name }}/bccvl_buildout
priority = 998
user = {{ user.name }}
stdout_logfile=syslog
stderr_logfile=syslog
autostart = true
numprocs = 1
startsecs = 10

; Wait for current tasks to finish at shutdown
; should be at least the time of the longest expected job
stopwaitsecs = 600

; When resorting to send SIGKILL, we want to kill all child processes as well
killasgroup = true
; Make sure subprocesses are stopped along with main process
stopasgroup = true
