# salt/roots/salt/dev/bccvl/data_mover_worker/data_mover_worker_supervisord.ini
[program:data_mover_worker]
command = scl enable python27 ". bin/activate; bin/celery worker --app=org.bccvl.tasks --loglevel=INFO  --queues=datamover --no-color --hostname=data_mover@%%h"
environment = CELERY_JSON_CONFIG="/home/data_mover/worker/celery.json"
process_name = data_mover_worker
directory = /home/data_mover/worker/
priority = 998
user = data_mover
stdout_logfile=syslog
stderr_logfile=syslog
autostart = true
numprocs = 1
startsecs = 10

; Wait for current tasks to finish at shutdown
; should be at least the time of the longest expected job
stopwaitsecs = 600

; When resorting to send SIGKILL, we want to kill all child processes as well
killasgroup = true
; Make sure subprocesses are stopped along with main process
stopasgroup = true