# salt/roots/salt/dev/bccvl/data_mover/data_mover_supervisord.ini
{% set user = salt['pillar.get']('data_mover:user', {'name': 'data_mover'}) %}

[program:data_mover]
command = scl enable python27 ". bin/activate; /home/{{ user.name }}/bccvl_data_mover/data_mover/bin/pserve data_mover.ini"
process_name = data_mover
directory = /home/{{ user.name }}/bccvl_data_mover/data_mover/
priority = 20
user = {{ user.name }}
stdout_logfile=syslog
stderr_logfile=syslog
autostart = true
numprocs = 1

; When resorting to send SIGKILL, we want to kill all child processes as well
killasgroup = true
; Make sure subprocesses are stopped along with main process
stopasgroup = true
