# salt/roots/salt/formulas/my-formulas/flower/flower_supervisord.ini
{% set cabundle = salt['pillar.get']('flower:cabundle', '/etc/pki/tls/cert.pem') %}
{% set auth = salt['pillar.get']('flower:auth') %}
{% if auth %}
{%   set auth = '--auth="' + auth + '"' %}
{% endif %}
; ==================================
; celery worker supervisor example
; ==================================
[program:flower]
; Set full path to celery program if using virtualenv
command=/home/flower/bin/flower --config="celeryconfig" {{ auth }}
environment=REQUESTS_CA_BUNDLE="{{ cabundle }}",PYTHONPATH="/home/flower/etc/config"
directory=/home/flower
user=flower
numprocs=1
stdout_logfile=syslog
stderr_logfile=syslog
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=998
