# salt/roots/salt/dev/bccvl/plone/plone_supervisord.ini
{% set user = salt['pillar.get']('plone:user', 'plon') %}
{% set plone_dir = '/home/' + user + '/bccvl_buildout' %}
{% set plone_bin = plone_dir + '/bin' %}

[program:zeo]
command = scl enable python27 ". ./bin/activate; {{ plone_bin }}/zeoserver fg"
process_name = zeo
directory = {{ plone_dir }}
priority = 10
user = plone
redirect_stderr = true
stdout_logfile = syslog
stderr_logfile = syslog
autostart = true
numprocs = 1
killasgroup = true
stopasgroup = true

{% set instances = salt['pillar.get']('plone:instances', {}) %}
{% for service in instances.keys() %}
[program:{{ service }}]
command = scl enable python27 ". ./bin/activate; {{ plone_bin }}/{{ service }} console"
process_name = {{ service }}
directory = {{ plone_dir }}
priority = 20
redirect_stderr = true
user = plone
{% if instances[service].get('debug', False) %}
autostart = false
{% endif %}
stdout_logfile = syslog
stderr_logfile = syslog
numprocs = 1
killasgroup = true
stopasgroup = true
{% endfor %}
