# salt/roots/salt/dev/bccvl/visualiser/visualiser_supervisord.ini
{% set user = salt['pillar.get']('visualiser:user', 'visualiser') %}

[program:visualiser]
command = scl enable python27 ". bin/activate; /home/{{ user }}/BCCVL_Visualiser/BCCVL_Visualiser/bin/pserve ploduction.ini"
process_name = visualiser
directory = /home/{{ user }}/BCCVL_Visualiser/BCCVL_Visualiser
priority = 20
redirect_stderr = true
user = {{ user }}
stdout_logfile=syslog
stderr_logfile=syslog
autostart = true
numprocs = 1

; When resorting to send SIGKILL, we want to kill all child processes as well
killasgroup = true
; Make sure subprocesses are stopped along with main process
stopasgroup = true